# infra/hermes-lite-vps/systemd/hermes-lite-maintenance.service.example
[Unit]
Description=Hermes Lite daily maintenance

[Service]
Type=oneshot
User=hermes
Group=hermes
Environment=HERMES_HOME=/home/hermes/.hermes-lite
ExecStart=/home/hermes/.hermes-lite/bin/backup-hermes-lite.sh
Nice=10
IOSchedulingClass=best-effort
