# infra/hermes-lite-vps/systemd/hermes-lite.service.example
[Unit]
Description=Hermes Lite Gateway
Documentation=https://github.com/NousResearch/hermes-agent
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=hermes
Group=hermes
WorkingDirectory=/home/hermes/.hermes-lite
EnvironmentFile=/home/hermes/.hermes-lite/.env
EnvironmentFile=/home/hermes/.hermes-lite/.env.lark
Environment=HERMES_HOME=/home/hermes/.hermes-lite
ExecStartPre=/usr/bin/test -x /home/hermes/.hermes-lite/venv/bin/hermes
ExecStart=/home/hermes/.hermes-lite/venv/bin/hermes gateway run
Restart=always
RestartSec=10s
TimeoutStopSec=120s
MemoryHigh=500M
MemoryMax=500M
TasksMax=64
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/hermes/.hermes-lite /home/hermes/.local /home/hermes/.cache
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
